# root/etc/services.d/sftpgo/run
#!/usr/bin/with-contenv bash

exec /usr/local/bin/sftpgo serve -c /config/sftpgo/