# root/etc/services.d/logs/run
#!/usr/bin/with-contenv bash

exec tail -f /var/log/sftpgo.log